# webpanel.service
[Unit]
Description=Web Panel Server Control
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/webpanel
Environment=PATH=/opt/webpanel/venv/bin
ExecStart=/opt/webpanel/venv/bin/python app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 